# systemd/letsencrypt.service
[Unit]
Description=Lets Encrypt Renew Certificate
After=network.target

[Service]
User=root
ExecStart=/root/certbot_renew/certbot_renew
